# deploy/systemd/hyacine-gallery-backend.service
[Unit]
Description=Hyacine Gallery backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hyacine-gallery
Group=hyacine-gallery
WorkingDirectory=/opt/hyacine-gallery/backend
EnvironmentFile=/etc/hyacine-gallery/hyacine-gallery.env
ExecStart=/bin/sh -lc 'exec /opt/hyacine-gallery/backend/.venv/bin/uvicorn app.main:app --host "${BACKEND_HOST:-127.0.0.1}" --port "${BACKEND_PORT:-8000}"'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
